# music-server@.service
[Unit]
Description=cubie-server Music Server (mDNS: cubie-server.local) for user %i
Documentation=https://github.com/snusnumrick/aiy_media_server
After=network.target network-online.target
Wants=network-online.target
Requires=network-online.target

[Service]
Type=simple

# Use %i (instance name) which becomes the username
# When installed as "music-server@bob.service", %i="bob"
User=%i
Group=%i

# Auto-detect home directory - %h/%i are available in all directives
WorkingDirectory=%h/music_server

# Auto-detect Python path
Environment=PATH=%h/music_server/music_server/bin
ExecStart=%h/music_server/music_server/bin/python %h/music_server/app.py

# Reload handler
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=music-server@%i

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

# Allow write access to music directory (auto-detected)
ReadWritePaths=%h/music_server/music

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
